# ansible/files/azure-pipelines-image-version-autogen.timer
[Unit]
Description=Update azure-pipelines-image-version routinely

[Timer]
OnActiveSec=1m
OnUnitActiveSec=1d

[Install]
WantedBy=timers.target
